# templates/supervisor/core.conf
[program:core-gunicorn]
command = {/project_dir/}gunicorn_script.bash
user = www-data                                                         
stdout_logfile = /path_to_logs/logs/gunicorn.log  
redirect_stderr = true                                    
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8